# roles/immich/templates/immich.service.j2
[Unit]
Description=Immich Media Server
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
User={{ ansible_user }}
Group={{ ansible_user }}
WorkingDirectory={{ podman_projects_dir }}/immich
ExecStart=/usr/bin/podman-compose up -d
ExecStop=/usr/bin/podman-compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
