# roles/disks/templates/fstrim.timer.j2
[Timer]
OnCalendar={{ ssd_trim_periodicity }}
